# com.gidiconnect.newsagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gidiconnect.newsagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/femimoritiwon/gidi-vibe-connect-1/scripts/lagos-news-agent.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/femimoritiwon/gidi-vibe-connect-1</string>

    <key>StartInterval</key>
    <integer>10800</integer>

    <key>StandardOutPath</key>
    <string>/Users/femimoritiwon/gidi-vibe-connect-1/logs/news-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/femimoritiwon/gidi-vibe-connect-1/logs/news-agent-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
